# deploy/flask_rest_client.ini
[program:flask_rest_client]
chdir=/home/bewi/flask_rest_client
home=/home/bewi/flask_rest_client
pythonpath=/home/bewi/flask_rest_client/ENV/bin/python
virtualenv=/home/bewi/flask_rest_client/ENV+
command=/bin/bash -c '~/flask_rest_client/ENV/bin/uwsgi uwsgi.ini'
directory=/home/bewi/flask_rest_client
autostart=yes
autorestart=yes
stopsignal=INT
stopasgroup=true
killasgroup=true

